# revisions/HAL/0.5.4-r1/60-HAL.rules
# Call HAL for all udev events
RUN+="/Programs/HAL/Current/libexec/hal.hotplug"

