# DBus/1.0.2/Resources/Tasks/MessageBus
#!/bin/sh
#
# messagebus:   The D-BUS systemwide message bus
#
# chkconfig: 345 97 03
# description:  This is a daemon which broadcasts notifications of system events \
#               and other messages. See http://www.freedesktop.org/software/dbus/
#
# processname: dbus-daemon
# pidfile: /Programs/DBus/Variable/run/dbus/pid
#

# Sanity checks.
#[ -x /Programs/DBus/0.95/bin/dbus-daemon ] || exit 0
mkdir -p /System/Variable/run/dbus/
rm /System/Variable/run/dbus/pid > /dev/null

# so we can rearrange this easily
RETVAL=0

/bin/dbus-daemon --system
RETVAL=$?
#sleep 2
exit $RETVAL
